# postgres/crontab
0 * * * * /usr/local/bin/backup_script.sh
